# pi/budget-eink.service
[Unit]
Description=Refresh budget summary on the Inky Impression
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=hobbes
EnvironmentFile=/etc/budget-eink/config.env
ExecStart=/home/hobbes/.virtualenvs/pimoroni/bin/python /home/hobbes/budget-eink/client.py
# The Inky 7.3" takes ~30s to refresh; give it room.
TimeoutStartSec=120
